# openwrt-ap-bedroom/overlay/etc/rc.local
if rm /etc/seccomp/umdns.json; then
  /etc/init.d/umdns restart || true
fi

# disable ethernet pause frames
ethtool -A eth0 autoneg off || :
ethtool -A lan1 autoneg off || :
ethtool -A lan2 autoneg off || :
ethtool -A lan3 autoneg off || :
ethtool -A lan4 autoneg off || :
