# openwrt-ap-balcony/overlay/etc/uci-defaults/61_batman_server
#!/bin/sh
uci set batman-adv.bat0.gw_mode='client'
uci commit batman-adv
