# openwrt-otg/overlay/etc/uci-defaults/21_system_settings
#!/bin/sh

uci -q batch <<-EOF >/dev/null
  set system.@system[0].hostname='otg-router'
  commit system
EOF

exit 0
