# openwrt-all/overlay/etc/uci-defaults/11_rootpwd
#!/bin/sh

source /etc/secrets/env

sed -i "s;^root.*;${root_shadow};" /etc/shadow
