# openwrt-all/overlay/etc/uci-defaults/99_rtreffer_authorized_keys
#!/bin/sh

cat >> /etc/dropbear/authorized_keys <<-_KEY_
ssh-rsa AAAAB3NzaC1yc2EAAAADAQABAAABAQCpQ7vS6xZGMUnjxSrl3d6Gdw7P1m+I7p9bQIk4UMWu3iIAWOfVlnU8GzqWWT7BFfpP2JMDx415yfrkw3dfY2zMdzT08DQKYy42eNGQ3CMwQ0v/NjjKCG4W6jTJux1PUjj6sUS9Z6ZjZ/ulknx9+5imfLzcTIsYX4RxLnk4P06UQpcr4KScxvnaAGO9uuBf+t/Lp9kFb00mv9/JU9KyXnIVj0vZRVSYK0t66UmUwmTVtPkSFIgKolDdRQKeLMEl1t0F7lwk4u/vV31jLw9aRX80yAtREE0Dfuw8/oDNPxdMEFJgLqW6o/KhjkvRTQc2PtPyRFRNmeaj0obbCpFsBfTR rtreffer@mactel
_KEY_
