# openwrt-all/overlay/etc/uci-defaults/99_rtreffer_authorized_keys
#!/bin/sh

cat >> /etc/dropbear/authorized_keys <<-_KEY_
ssh-rsa AAAAB3NzaC1yc2EAAAADAQABAAACAA1MNoLAjwxqd1+OpwpXcQf/k8GgtNmTklabU9QZkjCzN/XUPo8HnCY8ANWnwOkwZCigXMc6ur6idyt3LQ9HOx9Rap827k1aL6fZqHtP+pLSfkPtGgnbrf0IQQMhrXxxImdoWcFcCSXXWF/0IhspfkeiM+gRWLkUBGV526Su/oTJU8jTKBqW0koF+J5VQN/7H+ueqv8MRzyolL4BcgzWWL1Mz2TwspW1II7S0FQh1OkaGv8euKj2+XX/dIcxEBSNev3UlIE2Mtk7cCD4sKkqMgAJuiKw4ZRSSa7Zj95y9xq90Jt1LopkpDFNpqk4mppzHGON8KvDuVdevQ6Vte0ANcYa6w+lEj3Dn/kIK+QTSnrOa1LZcJa7wkgzKOw8cGInAJ3v2BFjAMNnjlvmx4Jtr81VrsOJenLf72pWgKUGNos1hgUBDyOGTeDmUPY4cacs0q/fA74NmmD6lVr3OQT7QCjVxf0Y8PwNDUZ4ZSZpsQUNQlkf1KD12501KYcv+lSudupCnaPAP2VDVXDkb/OQuuAag55FBX6YXUosiI69JsL0pD7tIoBE5FXKvZ8TA7Mc3FMOxr31eV22M7WsTBYlYGsPpIgLQ+hMcy3zGmRlbSG/VVim4flZhPfn9owwiOqfsUBBGtY2fG6RAAv6X0LX1r5TzS9P2bwDYeC51E4zHsvT rtreffer
_KEY_
