# openwrt-ap-living/overlay/etc/uci-defaults/61_batman_server
#!/bin/sh
uci set batman-adv.bat0.gw_mode='server'
uci commit batman-adv
