# openwrt-ap-living/overlay/etc/uci-defaults/61_wds
#!/bin/sh

uci -q batch <<-EOF >/dev/null
  set wireless.radio0_Strahlenkanone.wds=1
  commit wireless
EOF

exit 0
