# openwrt-ffgateway/overlay/etc/uci-defaults/98_umdns
#!/bin/sh

uci add_list 'umdms.@umdns[0].network'='wan'
uci commit umdns
