# openwrt-dsl/overlay/etc/rc.local
iptables -t raw -A PREROUTING -i dsl+ -j CT
iptables -t raw -A PREROUTING -i eth+ -j CT

exit 0
