# openwrt-ap-balcony/overlay/etc/rc.local
batctl routing_algo BATMAN_V || :
ip link add name bat0 type batadv || :
ip link set dev bat0 up || :
batctl if add wlan0 || :
batctl if add wlan1 || :
brctl addif br-lan bat0 || :

ip rule add fwmark 1 table ff.out
ip route add default via 192.168.2.1 table ff.out

ip route add 10.36.218.128/28 dev eth1.3
ip route add 10.0.0.0/8 via 10.36.218.129 dev eth1.3 src 10.36.218.139
