# openwrt-ap-max/overlay/etc/uci-defaults/99_firewall_off
../uci-defaults-available/99_firewall_off