# openwrt-ap-bedroom/overlay/etc/uci-defaults/21_system_settings
#!/bin/sh

uci -q batch <<-EOF >/dev/null
  set system.@system[0].hostname='ap-bedroom'
  commit system
EOF

uci -q batch <<-EOF >/dev/null
  set network.lan6=interface
  set network.lan6.proto='dhcpv6'
  set network.lan6.ifname='@lan'
  set network.lan6.reqprefix='no'
  commit network
EOF

uci -q batch <<-EOF >/dev/null
  set network.lan.ipaddr='192.168.2.5'
  del network.lan.ip6assign
  commit network
EOF

exit 0
