# openwrt-ap-bedroom/overlay/etc/uci-defaults/24_lldp
#!/bin/sh

uci -q batch <<-EOF >/dev/null
  add_list lldpd.config.interface=lan1
  add_list lldpd.config.interface=lan2
  add_list lldpd.config.interface=lan3
  add_list lldpd.config.interface=lan4
  commit lldpd
EOF

